# factoid-address-monitord@.service
[Unit]
Description=Monitor transactions for Factoid address %I
#Requires=factom-walletd.service
#After=factom-walletd.service

[Service]
Environment="NODE_PATH=/usr/lib/node_modules/"
EnvironmentFile=-/etc/default/factoid-address-monitord
ExecStart=/usr/bin/factoid-address-monitord --address %I --file /var/db/factoid-address-monitord/%I.csv $FACTOID_ADDRESS_MONITORD_OPTS
Restart=always

[Install]
WantedBy=default.target
